# templates/ledgersql.timer.erb
[Unit]
Description=Schedule ledger SQL sync

[Timer]
OnBootSec=<%= @bootdelay %>
OnUnitActiveSec=<%= @frequency %>

[Install]
WantedBy=multi-user.target
