# templates/ledgerdb.service.erb
[Unit]
Description=Update ledger database export
StartLimitBurst=2
StartLimitIntervalSec=271

[Service]
Type=simple
User=<%= @user %>
ExecStart=/usr/local/bin/ledgerdb sync <%= @datadir %>/config.yaml
Restart=on-failure
RestartSec=90
